# ops-bedrock/oneshot/s6-rc.d/op-init/up
#!/commands/with-contenv bash
/usr/local/bin/op-init.sh